# ansible/roles/hellard_generate_miner/templates/dashd-generate-miner.service.j2
[Unit]
After=network.target

[Service]
User={{ hellard_user }}
Group={{ hellard_group }}
ExecStart=/usr/local/bin/hellard-generate-miner.sh

[Install]
WantedBy=default.target
